# etc/rc.local
#!/bin/bash -e
#
# rc.local for raspimjpeg
#
# This script is executed at the end of each multiuser runlevel.
# value on error.
#
_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  printf "My IP address is %s\n" "$_IP"
fi

#sudo -u pi sh -c '/home/pi/software/motioneye/motioneye.py > /home/pi/software/motioneye/run/motioneye.log 2>&1' &


exit 0
