# config/a5-v11_15.05.1/etc/crontabs/root
#31 5 * * * reboot -f
*/5 * * * * /root/bin/manage3g.sh 
1 11 * * * /root/bin/manageWifi.sh start
#1 18 * * * /root/bin/manageWifi.sh stop     
